# zulip_botserver/zulip-botserver-supervisord.conf
[program:zulip-botserver]
command=zulip-botserver --config-file=<path/to/your/botserverrc> --hostname <address> --port <port>
startsecs=3
stdout_logfile=/var/log/zulip_botserver.log ; all output of your Botserver will be logged here
redirect_stderr=true
